# systemd/audio-profile-manager.service
[Unit]
Description=Fix HDMI audio after suspend/resume
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/audio-profile-manager
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
