# systemd/audio-profile-manager.service
[Unit]
Description=Fix HDMI audio after suspend/resume
After=graphical-session.target pipewire.service pipewire-pulse.service wireplumber.service

[Service]
ExecStart=%h/.local/bin/audio-profile-manager
Restart=on-failure

[Install]
WantedBy=graphical-session.target
